# rc.local
#!/bin/bash
#
# /etc/rc.local: Local multi-user startup script.
#

sudo -H -u eronarn screen -d -m -S irssi irssi
